# lib/hangctl/systemd/ablestack-vm-hangctl.service
[Unit]
Description=ABLESTACK VM hang controller (oneshot scan)
After=network.target libvirtd.service
Wants=libvirtd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ablestack_vm_hangctl scan
# scan ?ҙм—җ??dump/evidence/timeout???җмІҙ?ҒмңјлЎ?мІҳлҰ¬?ҳл?лЎ?systemd timeout?Җ ?үл„ү??
TimeoutStartSec=20min
# ?ӨнҢЁ ??systemdк°Җ ?¬мӢң?„н•ҳм§Җ ?ҠкІҢ (?¬мӢң?„лҠ” timerк°Җ ?ҙлӢ№)
Restart=no
